# CFG/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/fd/1
logfile_maxbytes=0

[program:ssh]
command=/usr/sbin/sshd -D
priority=100